# yocto/layers/meta-headunit/recipes-apps/mock-dbus/files/mock-theme.service
[Unit]
Description=Mock Theme D-Bus Service
Documentation=man:mock-theme(1)
After=dbus.service
Requires=dbus.service
Before=application-framework-manager.service

[Service]
Type=dbus
BusName=com.headunit.ThemeService
ExecStart=/usr/bin/mock-theme-service
Restart=always
RestartSec=3
TimeoutStartSec=10

# Environment
Environment="DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/0/bus"
Environment="XDG_RUNTIME_DIR=/run/user/0"

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mock-theme

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
