# yocto/layers/meta-headunit/recipes-apps/ivi-compositor/files/compositor.service
[Unit]
Description=HeadUnit IVI Nested Compositor (Wayland Client to Weston)
After=weston.service dbus.service dbus-session.service
Requires=weston.service dbus.service
Before=gearselector.service mediaplayer.service themecolor.service

[Service]
Type=simple
User=root
Environment="XDG_RUNTIME_DIR=/run/user/0"
Environment="WAYLAND_DISPLAY=wayland-0"
Environment="QT_QPA_PLATFORM=wayland"
Environment="QT_WAYLAND_DISABLE_WINDOWDECORATION=1"
Environment="IVI_SURFACE_ID=1001"
Environment="QT_LOGGING_RULES=*.debug=false;qt.qpa.*=true;qt.wayland.*=true"

ExecStartPre=/bin/mkdir -p /run/user/0
ExecStartPre=/bin/chmod 700 /run/user/0
ExecStartPre=/bin/sleep 2

ExecStart=/usr/bin/ivi-compositor

Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical.target