# yocto/layers/meta-headunit/recipes-core/systemd/boot-animation/boot-animation.service
[Unit]
Description=Boot Animation - LEXUS Logo
Before=weston.service compositor.service getty@tty1.service
After=sysinit.target
DefaultDependencies=no

[Service]
Type=forking
# Play video on framebuffer with loop
ExecStart=/bin/sh -c ' \
    /usr/bin/gst-launch-1.0 -q \
        filesrc location=/usr/share/boot-animation/lexus-logo.mp4 ! \
        decodebin ! videoconvert ! videoscale ! \
        video/x-raw,width=1024,height=600 ! \
        fbdevsink device=/dev/fb0 sync=false \
    & \
    echo $! > /run/boot-animation.pid \
'

# Kill the animation when Weston starts
ExecStop=/bin/sh -c 'kill $(cat /run/boot-animation.pid 2>/dev/null) 2>/dev/null || true'

# Clean up PID file
ExecStopPost=/bin/sh -c 'rm -f /run/boot-animation.pid'

RemainAfterExit=yes
StandardOutput=journal
StandardError=journal
TimeoutStartSec=10
TimeoutStopSec=2

[Install]
WantedBy=graphical.target
